# files/nixstatsagent.service
[Unit]
Description=Nixstatsagent

[Service]
ExecStart=/usr/local/bin/nixstatsagent
User=nixstats

[Install]
WantedBy=multi-user.target
